# etc/init/bat-worker.conf
description "Bat Worker"
author "longgeek <longgeek@thstack.com>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec start-stop-daemon --start --exec /usr/sbin/bat-worker
